# ops/systemd/arw-server-container.service
[Unit]
Description=ARW Unified Server (container)
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
EnvironmentFile=-/etc/default/arw-server
Environment=ARW_PORT=8091
Environment=ARW_DEBUG=0
Environment=ARW_BIND=127.0.0.1
Environment=ARW_IMAGE_OWNER=t3hw00t
Environment=ARW_IMAGE_TAG=latest
Environment=ARW_HTTP_TIMEOUT_SECS=20
Environment=ARW_HTTP_CONNECT_TIMEOUT_SECS=3
Environment=ARW_HTTP_TCP_KEEPALIVE_SECS=60
Environment=ARW_HTTP_POOL_IDLE_SECS=90
Environment=ARW_SAFE_MODE_ON_CRASH=1
User=%i
ExecStart=/usr/bin/docker run --rm \
  --name arw-server-%i \
  -p ${ARW_PORT}:${ARW_PORT} \
  -e ARW_PORT=${ARW_PORT} \
  -e ARW_BIND=${ARW_BIND} \
  -e ARW_DEBUG=${ARW_DEBUG} \
  -e ARW_ADMIN_TOKEN=${ARW_ADMIN_TOKEN} \
  ghcr.io/${ARW_IMAGE_OWNER}/arw-server:${ARW_IMAGE_TAG}
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
