# ops/systemd/arw-svc-native.service
[Unit]
Description=ARW Service (native)
After=network-online.target
Wants=network-online.target

[Service]
EnvironmentFile=-/etc/default/arw-svc
User=%i
ExecStart=/usr/local/bin/arw-svc
Restart=always
RestartSec=2
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ProtectControlGroups=true
ProtectClock=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallFilter=@system-service

[Install]
WantedBy=default.target

